# src/control-panel/src/main/shell/control-panel.service
[Unit]
Description=Control panel daemon
After=syslog.target
 
[Service]
ExecStart=/opt/control-panel/cpanel
 
[Install] 
WantedBy=multi-user.target